# host/roles/backup/templates/backup_service.timer.j2
[Unit]
Description=Triggers backup for {{ backup_service }}

[Timer]
Unit=backup@{{ backup_service }}.service
OnCalendar=*-*-* 05:00:00

[Install]
WantedBy=timers.target
