# ansible/roles/services/ovh_dyndns_updater/templates/ovh_dyndns_updater.timer
[Unit]
Description=Updates OVH DynDNS record every 10min

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
