# ansible/roles/services/vault/templates/unseal-vault-check.timer
[Unit]
Description=Check if vault needs unsealing

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
