# host_system/media-data2.automount
[Unit]
Description=Automount for second data drive

[Automount]
Where=/media/data2/
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
