# ansible/roles/backup/templates/backup-remote-service.timer.j2
[Unit]
Description=Triggers weekly remote backup

[Timer]
Unit=ansible-backup.service
OnCalendar=Mon 05:00

[Install]
WantedBy=timers.target
