# ansible/roles/backup/templates/backup-remote.service.j2
[Unit]
Description=Perform weekly remote backups

[Service]
Type=oneshot
ExecStart={{ homeserver_repo_dir }}/ansible/backup.sh remote all
