# ansible/roles/backup/templates/backup.service.j2
[Unit]
Description=Perform local backups

[Service]
Type=oneshot
ExecStart={{ homeserver_repo_dir }}/ansible/backup.sh local all
