# host/roles/server/templates/second_drive.automount.j2
[Unit]
Description=Automount for second data drive

[Automount]
Where={{ drives.backup.mount }}
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
