# ansible/roles/server/templates/remote_backup.automount.j2
[Unit]
Description=Automount for remote backup drive

[Automount]
Where={{ drives.remote_backup.mount }}
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
